# deploy/csbg-intake.service
[Unit]
Description=CSBG Client Intake System (Next.js on :3100)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dchevalier
Group=dchevalier
WorkingDirectory=/opt/csbg-intake
EnvironmentFile=/etc/csbg-intake.env
ExecStart=/usr/bin/node node_modules/next/dist/bin/next start -p 3100
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
